# config/autostart/GoXLR_daemon.desktop
[Desktop Entry]
Exec=/usr/bin/goxlr-daemon
Icon=application-x-executable
Name=goxlr-daemon
Type=Application
X-KDE-AutostartScript=true