# config/autostart/GoXLR_daemon.desktop
[Desktop Entry]
Exec=goxlr-daemon
Icon=application-x-shellscript
Name=goxlr-daemon
Type=Application
X-KDE-AutostartScript=true